# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null

[program:anvil]
command=/usr/local/bin/anvil -f 'https://mainnet.base.org' --host 0.0.0.0 -m 'talent catch cousin flavor first giant wisdom accident soup engine corn taste' -a 100
logfile=/dev/stdout
user=root
autostart=true
priority=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:healthcheck]
command=/bin/sh -c "while true; do curl localhost:8545; if [ $? -eq 0 ]; then supervisorctl start deploy; break; else sleep 1; fi; done"
autostart=true
autorestart=false
startsecs=0

[program:deploy]
command=forge script script/Deploy.s.sol --rpc-url http://127.0.0.1:8545 --broadcast -vvvv
user=root
autostart=false
autorestart=false
priority=2
environment=PRIVATE_KEY=0x601f6fa6e4b65cf41012631a3e1062f42a8e5a3bc78c4972f95e4412f5be2b3d
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
